# python-learning-site.service
[Unit]
Description=Python Learning Site
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/.openclaw/workspace/python-learning-site
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 优雅关闭超时时间
TimeoutStopSec=30

# 防止快速重启循环
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
